# lib/deploy/ansible/roles/common/files/upstart-unicorn.conf
description "Recruiter"
author "Jeffrey Horn <jrhorn424@gmail.com>"

env APP_ROOT=/u/apps/recruiter
env APP_PATH=$APP_ROOT/current

start on virtual-filesystems
stop on runlevel [06]

env PATH=/opt/ruby/shims:/opt/ruby/rbenv/bin:/usr/local/bin:/usr/bin:/bin

env RAILS_ENV=production

chdir $APP_PATH

pre-start exec bundle exec unicorn -D -c config/unicorn.rb --E production

post-stop exec kill `cat $APP_PATH/pids/unicorn.pid`
